# rtsp2mjpg/rootfs/etc/services.d/ffmpeg/run
#!/usr/bin/with-contenv bashio

# Configuratie ophalen
rtsp_url=$(bashio::config 'rtsp_url')
ffmpeg_input_opts=$(bashio::config 'ffmpeg_input_opts')
ffmpeg_output_opts=$(bashio::config 'ffmpeg_output_opts')
ffserver_log_level=$(bashio::config 'ffserver_log_level')
ffmpeg_log_level=$(bashio::config 'ffmpeg_log_level')

# Validatie
if [[ -z "${rtsp_url}" ]]; then
    bashio::exit.nok "RTSP URL is required"
fi

# Omgevingsvariabelen exporteren
export RTSP_URL="${rtsp_url}"
export FFMPEG_INPUT_OPTS="${ffmpeg_input_opts}"
export FFMPEG_OUTPUT_OPTS="${ffmpeg_output_opts}"
export FFSERVER_LOG_LEVEL="${ffserver_log_level}"
export FFMPEG_LOG_LEVEL="${ffmpeg_log_level}"

bashio::log.info "Starting RTSP2MJPG addon..."
bashio::log.info "RTSP URL: ${rtsp_url}"
